# Cron Docker/crontab
* * * * * python3 /home/generate.py && mv /root/index.html /home/

